# ami/files/hazelcast.timer
[Unit]
Description=Hazelcast server

[Timer]
OnBootSec=30

[Install]
WantedBy=multi-user.target
